# skull-king.service
[Unit]
Description=Skull King Game Server (PM2)
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=sydnec
Group=sydnec
WorkingDirectory=/home/sydnec/skullking_js
ExecStart=/usr/bin/pm2 start ecosystem.config.js --env production
ExecStop=/usr/bin/pm2 stop skullking
ExecReload=/usr/bin/pm2 restart skullking
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=-/home/sydnec/skullking_js/.env.production
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skull-king

# Limites de sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/sydnec/skullking_js
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Limites de ressources
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
